# config/launchd/com.paperpile-obsidian.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.paperpile-obsidian</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/mtsr95/.bun/bin/bun</string>
    <string>run</string>
    <string>/Users/mtsr95/Documents/GitHub/paperpile-obsidian/src/index.ts</string>
    <string>daemon</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/mtsr95</string>
    <key>PATH</key>
    <string>/Users/mtsr95/.bun/bin:/usr/local/bin:/usr/bin:/bin:/usr/local/sbin</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/mtsr95/Documents/GitHub/paperpile-obsidian</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/mtsr95/.paperpile-obsidian/daemon.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/mtsr95/.paperpile-obsidian/daemon.error.log</string>

  <!-- Restart after 30s if it crashes -->
  <key>ThrottleInterval</key>
  <integer>30</integer>
</dict>
</plist>
